# supervisor.conf
[program:inspectify]
process_name=%(program_name)s_%(process_num)02d
command=/home/inspectify/app/venv/bin/gunicorn --workers 4 --bind 127.0.0.1:8000 --config /home/inspectify/app/gunicorn_config.py inspectify.wsgi:application
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/home/inspectify/logs/gunicorn_stdout.log
stderr_logfile=/home/inspectify/logs/gunicorn_stderr.log
user=inspectify
directory=/home/inspectify/app
environment=PATH="/home/inspectify/app/venv/bin",DB_PASSWORD="your-secure-password"

[program:celery-worker]
process_name=%(program_name)s_%(process_num)02d
command=/home/inspectify/app/venv/bin/celery -A inspectify worker -l info
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/home/inspectify/logs/celery_stdout.log
stderr_logfile=/home/inspectify/logs/celery_stderr.log
user=inspectify
directory=/home/inspectify/app

[program:celery-beat]
process_name=%(program_name)s
command=/home/inspectify/app/venv/bin/celery -A inspectify beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/home/inspectify/logs/celery_beat_stdout.log
stderr_logfile=/home/inspectify/logs/celery_beat_stderr.log
user=inspectify
directory=/home/inspectify/app

[group:inspectify]
programs=inspectify,celery-worker,celery-beat
priority=999
